# ansible/storage/roles/apps.storage/files/docker-cleanup.timer
[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitInactiveSec=12h

[Install]
WantedBy=timers.target
