# ansible/kvm/roles/node-exporter.kvm/templates/node-exporter.service.j2
[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
User=root
ExecStart=/usr/local/bin/node-exporter

[Install]
WantedBy=multi-user.target
