# ansible/kubernetes/playbooks/templates/etcd-defrag.timer.j2
[Unit]
Description=Run etcd-defrag.service every day
After=k3s.service
[Timer]
OnCalendar=*-*-* 01:00:0
[Install]
WantedBy=multi-user.target
