# ansible/kubernetes/playbooks/files/stale-containers.service
[Unit]
Description=Stale containers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s crictl rmi --prune
