# ansible/pikvm/roles/promtail.pikvm/templates/promtail.service.j2
[Unit]
Description=Promtail

[Service]
TimeoutStartSec=0
User=root
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/promtail.yaml

[Install]
WantedBy=multi-user.target
